# beemaster_connector.service
[Unit]
Description=Beemaster Connector
Author=Beemaster
After=docker.service

[Service]
Restart=always
WorkingDirectory=/home/beemaster/mp-ids/mp-ids-hp/connector/src
ExecStart=/usr/bin/sh /home/beemaster/mp-ids/mp-ids-hp/picontrol.sh conn start
ExecStop=/usr/bin/sh /home/beemaster/mp-ids/mp-ids-hp/picontrol.sh conn stop
